# catalogs/sigma/ast/c5c5f538-8d00-5a54-84ab-cda515fc55c1.json
{
  "id": "c5c5f538-8d00-5a54-84ab-cda515fc55c1",
  "catalog": "sigma",
  "name": "Suspicious Registry Run Key Modification",
  "description": "Detects modification of Run/RunOnce registry keys for persistence",
  "severity": "medium",
  "mitre_techniques": [
    "attack.t1547.001"
  ],
  "event_categories": [
    "registry"
  ],
  "conditions": [
    {
      "field": "registry.path",
      "raw_field": "TargetObject",
      "operator": "like~",
      "values": [
        "*\\CurrentVersion\\Run*"
      ],
      "raw_values": [
        "*\\CurrentVersion\\Run*"
      ]
    }
  ],
  "raw_query": "registry where registry.path like~ \"*\\\\CurrentVersion\\\\Run*\"",
  "language": "sigma",
  "translated_query": "registry where registry.path like~ \"*\\\\CurrentVersion\\\\Run*\"",
  "source_path": "demo",
  "metadata": {}
}