# catalogs/elastic/ast/3bb014f3-8063-5dfa-81b3-8a2c6ff7acde.json
{
  "id": "3bb014f3-8063-5dfa-81b3-8a2c6ff7acde",
  "catalog": "elastic",
  "name": "Windows Registry Run Key Persistence",
  "description": "Detects addition of entries to registry Run keys",
  "severity": "medium",
  "mitre_techniques": [
    "attack.t1547.001"
  ],
  "event_categories": [
    "registry"
  ],
  "conditions": [
    {
      "field": "registry.path",
      "raw_field": "registry.path",
      "operator": "like~",
      "values": [
        "*\\\\CurrentVersion\\\\Run*"
      ],
      "raw_values": [
        "*\\\\CurrentVersion\\\\Run*"
      ]
    }
  ],
  "raw_query": "registry where registry.path like~ \"*\\\\CurrentVersion\\\\Run*\"",
  "language": "eql",
  "translated_query": null,
  "source_path": "demo",
  "metadata": {}
}